# systemd/fermentation-dashboard.service
[Unit]
Description=Fermentation Web Dashboard
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/fermentation-pump-controller
ExecStart=/usr/bin/python3 /home/pi/fermentation-pump-controller/src/web_dashboard.py
Restart=always

[Install]
WantedBy=multi-user.target
