# systemd/pump-monitoring.service
[Unit]
Description=Fermentation Pump Morning Cycle
After=network.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/fermentation-pump-controller
ExecStart=/usr/bin/python3 /home/pi/fermentation-pump-controller/src/pump_control.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
